# init/beyond-sync-api.service
[Unit]
Description=Beyond Sync API Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=beyond-sync
Group=beyond-sync
WorkingDirectory=/opt/beyond-sync
EnvironmentFile=/opt/beyond-sync/.env
ExecStart=/opt/beyond-sync/beyond-sync-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target